# tailscale/rootfs/etc/s6-overlay/s6-rc.d/addon-services/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Configures and registers addon services with Tailscale
# ==============================================================================

if ! bashio::config.has_value "addon_services"; then
    bashio::log.info "No addon services configured"
    touch /run/addon_services_ready
    exit 0
fi

bashio::log.info "Setting up addon services..."

# Iterate through each addon service
bashio::config "addon_services" | jq -c '.[]' | while read -r service; do
    port=$(echo "$service" | jq -r '.port // empty')
    service_name=$(echo "$service" | jq -r '.service_name // empty')
    host=$(echo "$service" | jq -r '.host // "localhost"')
    funnel=$(echo "$service" | jq -r '.funnel // false')
    
    if [ -z "$port" ] || [ -z "$service_name" ]; then
        bashio::log.error "Invalid service config - port: '$port', service_name: '$service_name'"
        bashio::log.error "Full config: $service"
        continue
    fi
    
    # Build the target URL
    target_url="http://$host:$port"
    
    bashio::log.info "Registering service: $service_name (target: $target_url, funnel: $funnel)"
    
    if [ "$funnel" = "true" ]; then
        bashio::log.debug "Executing: tailscale serve funnel tcp:443/$service_name $target_url"
        if tailscale serve funnel "tcp:443/$service_name" "$target_url"; then
            bashio::log.info "Successfully set up funnel for $service_name"
        else
            local exit_code=$?
            bashio::log.error "Failed to set up funnel for $service_name (exit code: $exit_code)"
        fi
    else
        bashio::log.debug "Executing: tailscale serve svc:$service_name $target_url"
        if tailscale serve "svc:$service_name" "$target_url"; then
            bashio::log.info "Successfully set up serve for $service_name"
        else
            local exit_code=$?
            bashio::log.error "Failed to set up serve for $service_name (exit code: $exit_code)"
        fi
    fi
done

bashio::log.info "Addon services setup complete"
touch /run/addon_services_ready
