# tailscale/rootfs/etc/s6-overlay/s6-rc.d/addon-services/up
/etc/s6-overlay/s6-rc.d/addon-services/run