# tailscale/rootfs/etc/s6-overlay/s6-rc.d/addon-services/up
#!/command/execlineb -P
s6-test -f /run/addon_services_ready
